# files/usr/lib/systemd/user/hide-waydroid-apps.path
[Unit]
Description=Watch for new Waydroid apps

[Path]
# Watch the user's application folder for changes
PathChanged=%h/.local/share/applications

[Install]
WantedBy=default.target