# files/usr/lib/systemd/user/hide-waydroid-apps.service
[Unit]
Description=Hide Waydroid Apps from App Grid
After=gnome-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/hide-waydroid-apps