# files/usr/lib/systemd/system/firefox-prefs.service
[Unit]
Description=Configure firefox preferences
After=flatpak-preinstall.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/firefox-prefs
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target