# files/usr/lib/systemd/system/libvirt-backup.service
[Unit]
Description=Backup libvirt VM XML configurations
After=network.target

[Service]
Type=oneshot
User=root
ExecStartPre=/bin/mkdir -p /home/jonah/vms
ExecStart=/bin/bash -c 'for vm in $(virsh list --all --name); do [ -n "$vm" ] && virsh dumpxml "$vm" > "/home/jonah/vms/${vm}.xml"; done'