# config/files/usr/lib/systemd/system/update.timer
[Unit]
Description=Update image timer
Requires=update.service

[Timer]
Unit=update.service
OnCalendar=daily

[Install]
WantedBy=timers.target