# config/files/usr/lib/systemd/system/update.service
[Unit]
Description=Update image everyday

[Service]
Type=oneshot
ExecStart=ujust update


[Install]
WantedBy=multi-user.target