# config/files/usr/lib/systemd/system/update.timer
[Unit]
Description=Update image timer
Requires=update.service

[Timer]
Unit=update.service
OnUnitInactiveSec=24h

[Install]
WantedBy=timers.target